# old/dtsx/templates/dtsx.ini.erb
[program:DTSX<%= name %>]
command=/usr/bin/wine 'c:\servers\jdk1.6\bin\java.exe' -cp 'z:\usr\adagio\Adagio_Lib-<%= adagio_version %>-client-full.jar' -DJINTEGRA_NATIVE_MODE '-Ddtsx.server.configfile=c:\servers\DTSX\dtsx-server.config' adk.dts.server.DTSXAdminServer
; Use this flag to set the Apartment threading model: -DJINTEGRA_COINIT_VALUE=2
directory=/usr/indesign        ; directory to cwd to before exec (def no cwd)
priority=150                   ; the relative start priority (default 999)
startsecs=30                   ; number of secs prog must stay running (def. 1)
user=indesign                     ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
environment=DISPLAY=':<%= name %>',WINEPREFIX='/usr/indesign/dtsx<%= name %>',WINEDEBUG='err-ole'
