# old/abc/templates/ReceiptProcessServer.ini.erb
[program:ReceiptProcessServer]
command=/usr/java/default/bin/java -cp /usr/adagio/Adagio_Lib-<%= adagio_version %>-jar-with-dependencies.jar:/usr/pts/PTS-<%= pts_version %>.jar:/usr/pts/antlr-2.7.6.jar:/usr/pts/commons-codec-1.2.jar:/usr/pts/commons-collections-3.1.jar:/usr/pts/commons-httpclient-3.1.jar:/usr/pts/commons-logging-1.0.4.jar:/usr/pts/db2java-2.10.jar:/usr/pts/db2jcc-2.10.jar:/usr/pts/dom4j-1.6.1.jar:/usr/pts/ejb3-persistence-1.0.2.GA.jar:/usr/pts/hibernate-3.3.1.GA.pom:/usr/pts/hibernate-annotations-3.4.0.GA.jar:/usr/pts/hibernate-commons-annotations-3.1.0.GA.jar:/usr/pts/hibernate-core-3.3.0.SP1.jar:/usr/pts/javassist-3.4.ga.jar:/usr/pts/jta-1.1.jar:/usr/pts/log4j-1.2.11.jar:/usr/pts/slf4j-api-1.4.2.jar:/usr/pts/slf4j-api-1.5.2.jar:/usr/pts/slf4j-simple-1.5.0.jar:/usr/pts/xml-apis-1.0.b2.jar:/usr/pts/xmlpublic-2.1.0.jar:/usr/pts/xpp3-min-1.1.4c.jar:/usr/pts/xstream-1.3.jar:/usr/adagio/db2jcc_license_cu.jar -Dadminserver.serverid=<%= serverid %> -Dadminserver.configfile=/usr/adagio/jndi.properties adagio.server.receiptprocess.ReceiptProcessAdminServer
;process_name=%(program_name)s ; process_name expr (default %(program_name)s)
;numprocs=1;          ; number of processes copies to start (def 1)
;directory=/usr/indesign        ; directory to cwd to before exec (def no cwd)
;umask=022                     ; umask for process (default None)
priority=400                   ; the relative start priority (default 999)
;autostart=true                ; start at supervisord start (default: true)
;autorestart=true              ; retstart at unexpected quit (default: true)
startsecs=30                   ; number of secs prog must stay running (def. 1)
;startretries=3                ; max # of serial start failures (default 3)
;exitcodes=0,2                 ; 'expected' exit codes for process (default 0,2)
;stopsignal=INT                 ; signal used to kill process (default TERM)
;stopwaitsecs=120               ; max num secs to wait b4 SIGKILL (default 10)
user=nobody                     ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
;environment=DISPLAY=':0'
